# contrib/ezweb.service
[Unit]
Description=EzWeb Site Manager
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/opt/ezweb
EnvironmentFile=/opt/ezweb/.env
ExecStart=/opt/ezweb/ezweb
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
